# systemd/sync-passwords.service
[Unit]
Description=Synchronise the passwords with the cloud

[Service]
Type=oneshot
ExecStart=/home/hakmad/.bin/sync-passwords.sh
